# Expt/setting/nks/skel/.config/autostart/firstboot.desktop
[Desktop Entry]
Name=SetupLive
Comment=SetupLive
Exec=/usr/bin/first-boot-check
StartupNotify=false
Terminal=false
Type=Application
X-KDE-PluginInfo-Author=SasaKen
X-KDE-PluginInfo-Email=kenta@nksteam.com
